# supervisord.conf.sample
[program:solmanager]
directory=#SOL_MANAGER_PATH
command=screen -LmDS solmanager /usr/bin/python solmanager.py
autostart=true
autorestart=true
stderr_logfile=sol.err.log
stdout_logfile=sol.out.log
